# srv/salt/supervisor/supervisord.conf
[program:realms]
command=/home/deploy/virtualenvs/realms/bin/python /home/deploy/realms/manage.py