# srv/salt/supervisor/supervisord.conf
[program:realms]
user=deploy
command=/home/deploy/realms/virtualenvs/realms/bin/python /home/deploy/realms/manage.py server